# supervisor.conf
[program:app]
command = /home/frank/ak/supervisor/watermark-site/env/bin/python manage.py db init
directory = /home/frank/ak/supervisor/
